# base/etc/cron.d/osg-configure
PATH=/usr/sbin:/usr/bin:/bin:$PATH
*/15 * * * * root /usr/local/bin/osg-configure-cron.sh
